# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --socket :8001 --wsgi-file /code/composeexample/wsgi.py

[program:nginx-app]
command = /usr/sbin/nginx
